# src/etc/services-config/supervisor/supervisord.d/redis-server-bootstrap.conf
[program:redis-server-bootstrap]
priority = 6
command = /usr/sbin/redis-server-bootstrap
autostart = %(ENV_REDIS_AUTOSTART_REDIS_BOOTSTRAP)s
startsecs = 0
startretries = 0
autorestart = false
redirect_stderr = true
stdout_logfile = /var/log/redis-server-bootstrap
stdout_events_enabled = true